# openrc
#!/sbin/openrc-run

name="short-url"
description="URL Shortener Service"

command="/usr/bin/node /opt/short-url/server.js"
command_user="shorturl"
directory="/opt/short-url"

# Environment variables
export PORT=3000
export BASE_URL=https://yourdomain.com
export DB_PATH=/var/lib/short-url/urls.db

# Service dependencies
depend() {
    need net
    after firewall
}

# Pre-start checks
start_pre() {
    checkpath -d -o shorturl:shorturl /opt/short-url
    checkpath -d -o shorturl:shorturl /var/lib/short-url
}
