# openrc
#!/sbin/openrc-run

name="short-url"
description="URL Shortener Service"

command="/usr/bin/node"
command_args="/opt/short-url/server.js"
command_user="shorturl:shorturl"
directory="/opt/short-url"

command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

# Logging configuration
stdout_logger="exec logger -t $RC_SVCNAME"
stderr_logger="exec logger -t $RC_SVCNAME"

# Environment variables
export PORT=3000
export BASE_URL=https://yourdomain.com
export DB_PATH=/var/lib/short-url/urls.db

# Service dependencies
depend() {
    need net
    after firewall
}

# Pre-start checks
start_pre() {
    checkpath -d -o shorturl:shorturl /run
    checkpath -d -o shorturl:shorturl /var/lib/short-url
}
